# packer/scripts/node-server.service
[Unit]
Description=my application
Documentation=https://example.com
After=network.target

[Service]

Type=simple
User=csye6225
Environment=DB_USERNAME=webapp
Environment=DB_PASSWORD=helloworld
Environment=DB_NAME=webapp
Environment=DB_HOST=127.0.0.1
ExecStart=/usr/bin/node /home/csye6225/webapp/server.js
Restart=on-failure

[Install]
WantedBy=multi-user.target